# docker-files/etc/supervisor/conf.d/cron.conf
[program:cron]
command=/etc/init.d/cron start
user=root
autostart=true
stdout_logfile=/var/log/cron.log
